# etc/profile.d/i18n.sh
# Set up i18n variables

export LANG=en_GB.ISO-8859-1
