# etc/profile.d/rustc.sh
# Begin /etc/profile.d/rustc.sh
pathprepend /opt/rustc/bin   PATH
# End /etc/profile.d/rustc.sh
